# docker/tomcat/tomcat.ini
[program:tomcat]
command=/opt/tomcat.sh
autorestart=true
stdout_logfile=/var/log/catalina.log
