# docker/inductor/inductor.ini
[program:inductor]
command=/opt/inductor.sh
autorestart=false
startretries=0
stdout_logfile=/var/log/inductor.log
