# docker/postgres/postgres.ini
[program:postgres]
command=/usr/pgsql-9.2/bin/postgres -D /var/lib/pgsql/9.2/data -i -h postgres
user=postgres
autorestart=true

[program:db]
autorestart=false
startretries=0
redirect_stderr=true
stdout_events_enabled=1
stdout_logfile=/var/log/db.log
command=/etc/supervisord.d/db.sh
