# docker/postgres/postgres.ini
[program:postgres]
command=/usr/pgsql-9.2/bin/postgres -D /var/lib/pgsql/9.2/data -i -h postgres
user=postgres
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:db]
command=/etc/supervisord.d/db.sh
autorestart=false
startretries=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
