# docker/rails/display.ini
[program:display]
command=/opt/display.sh
autorestart=true
stdout_logfile=/var/log/display.log
