# docker/elasticsearch/elasticsearch.ini
[program:elasticsearch]
command=/opt/elasticsearch.sh
stdout_logfile=/var/log/elasticsearch.log

[program:search]
command=/opt/search.sh
autorestart=true
redirect_stderr=true
stdout_events_enabled=1
stdout_logfile=/var/log/search.log
