# docker/logstash/logstash.ini
[program:logstash]
command=/opt/logstash.sh
stdout_logfile=/var/log/logstash.log
