# docker/cassandra/cassandra.ini
[program:cassandra]
command=/opt/cassandra/bin/cassandra -f
stdout_logfile=/var/log/cassandra.log
autorestart=true
