# docker/activemq/activemq.ini
[program:activemq]
command=/opt/activemq.sh
autorestart=true
